# workspace/projects/new-kanban-setup/.systemd-user/mission-control.target
[Unit]
Description=Mission Control runtime target
Wants=mission-control-api.service mission-control-ui.service
After=network-online.target
